# vscode/rootfs/etc/s6-overlay/s6-rc.d/init-code-server/run
#! /command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Studio Code Server
# Sets up code-server.
# ==============================================================================
declare config_path

if bashio::config.has_value 'config_path'; then
  config_path=$(bashio::config 'config_path')
  if ! bashio::fs.directory_exists "${config_path}"; then
    bashio::exit.nok "Configured config path does not exists"
  fi
fi

# Ensure persistent data folder exists.
if ! bashio::fs.directory_exists '/data/vscode'; then
  mkdir -p /data/vscode/extensions \
    || bashio::exit.nok "Could not create persistent storage folder."
fi

# Clean up copies of extensions we deliver from the persistent storage
while read -r ext; do
  extension="${ext%%@*}"
  # shellcheck disable=SC2086
  rm -f -r /data/vscode/extensions/${extension,,}*
done < /root/vscode.extensions
bashio::log.info "removed extensions that are provided in persistent storage";

# Ensure user extensions folder exists
mkdir -p /data/vscode/extensions

# Sets up default machine settings on first start.
mkdir -p /data/vscode/Machine \
  || bashio::exit.nok "Could not create persistent storage folder."
cp /root/.code-server/settings.json /data/vscode/Machine/settings.json
bashio::log.info "setup machine settings";
