# vscode/rootfs/etc/s6-overlay/s6-rc.d/init-code-server/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Studio Code Server
# Sets up code-server.
# ==============================================================================
declare config_path

# List of previous config hashes, to allow upgrade "default" configs.
readonly -a PREVIOUS_DEFAULT_CONFIG_HASHES=(
  0d05244fb32cc1991b456d56c71d303e0c3303f2ecb29c0329b9373799263efbc135d298c2b7d5e3035f5f3bf7eaf03fe83d12620146f1e31c49e396de81dc48
)

if bashio::config.has_value 'config_path'; then
  config_path=$(bashio::config 'config_path')
  if ! bashio::fs.directory_exists "${config_path}"; then
    bashio::exit.nok "Configured config path does not exists"
  fi
fi

# Ensure persistent data folder exists.
if ! bashio::fs.directory_exists '/data/vscode'; then
  mkdir -p /data/vscode/extensions \
    || bashio::exit.nok "Could not create persistent storage folder."
fi

# Clean up copies of extensions we deliver from the persistent storage
while read -r ext; do
  extension="${ext%%@*}"
  # shellcheck disable=SC2086
  rm -f -r /data/vscode/extensions/${extension,,}*
done < /root/vscode.extensions

# Ensure user extensions folder exists
mkdir -p /data/vscode/extensions

# Sets up default user settings on first start.
if ! bashio::fs.file_exists '/data/vscode/User/settings.json'; then
  mkdir -p /data/vscode/User \
    || bashio::exit.nok "Could not create persistent storage folder."

  cp /root/.code-server/settings.json /data/vscode/User/settings.json
fi

# Upgrade settings.json is still default from previous version.
current=$(sha512sum /data/vscode/User/settings.json|cut -d " " -f 1)
if [[ " ${PREVIOUS_DEFAULT_CONFIG_HASHES[*]} " == *" ${current} "* ]]; then
  cp /root/.code-server/settings.json /data/vscode/User/settings.json
fi
